# build/crontab
# m h dom mon dow command

# Clean out the cache as it gets very full
39 * * * *    /usr/sbin/tmpwatch 12 /home/francis/pwcache


