# build/crontab
# m h dom mon dow command

27 1 * * *       /home/francis/bin/recent-check
39 * * * *    /usr/sbin/tmpwatch 12 /home/francis/pwcache


