# build/remotescripts/crontab
# m h dom mon dow command

27 1 * * *       /home/francis/bin/recent-check
39 2,14 * * *    /usr/sbin/tmpwatch 24 /home/francis/pwcache

